# docker/cron/crontab
*/2 * * * * php /var/www/scripts/email.php >> /var/log/cron/cron.loh 2>&1